# sample_osquery_checks/Windows/jan_ai_local_llm_detection.yml
title: Jan.ai Local LLM Detection
id: 2f3a4b5c6d7e8f9a0b1c2d3e4f5a6b7c
description: |
    Detects the presence of Jan.ai, an open-source desktop application for running large
    language models entirely on-device. Jan.ai downloads and stores model weight files in
    %USERPROFILE%\.jan\models\ and provides an OpenAI-compatible API server on port 1337
    (the "Jan API Server"). The local server can be accessed by any application on the device
    or, if the user changes the bind address, from the local network. As with other local LLM
    servers (Ollama, LM Studio), all inference happens locally so outbound DLP controls have
    no visibility into data processed by Jan.ai. Jan.ai additionally supports MCP (Model
    Context Protocol) tool integrations that can grant local AI models access to the
    filesystem, shell, and other device resources. Organizations should classify Jan.ai under
    the same risk category as Ollama and LM Studio. This detection evaluates startup items,
    file paths, running processes, installed programs, listening ports, and prefetch entries.
    A score greater than 1 (at least two positive indicators) is required to confirm detection.
references:
    - https://jan.ai
    - https://www.okta.com/blog/threat-intelligence/detecting-openclaw-advanced-posture-checks/
author:
    - rafa.bono@okta.com
platform: Windows
query: |
    WITH startup_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM startup_items
        WHERE name LIKE '%jan%'
    ),
    file_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM file
        WHERE path LIKE '%\Users\%%\AppData\Local\Programs\Jan\Jan.exe'
            OR path LIKE '%\Users\%%\.jan\models\%'
            OR path LIKE '%\Users\%%\.jan\settings.json'
            OR path LIKE '%\Users\%%\AppData\Roaming\Jan\%'
    ),
    process_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM processes
        WHERE name = 'Jan.exe'
            OR name LIKE '%jan-ai%'
    ),
    programs_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM programs
        WHERE name LIKE '%Jan%'
            AND name NOT LIKE '%January%'
    ),
    netports_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM listening_ports
        WHERE port = '1337'
            AND path LIKE '%Jan%'
    ),
    prefetch_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM prefetch
        WHERE path LIKE '%\Jan\Jan.exe%'
            OR path LIKE '%\jan-ai%'
    ),
    sockets_jan AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM process_open_sockets
        WHERE path LIKE '%\Jan\Jan.exe'
    ),
    final_score AS (
        SELECT
            startup_jan.total
            + file_jan.total
            + process_jan.total
            + programs_jan.total
            + netports_jan.total
            + prefetch_jan.total
            + sockets_jan.total
            AS score
        FROM startup_jan, file_jan, process_jan, programs_jan,
             netports_jan, prefetch_jan, sockets_jan
    )
    SELECT
        CASE
            WHEN score <= 1 THEN 0
            WHEN score > 1 THEN 1
        END AS jan_ai_detected
    FROM final_score
    ;
