# sample_osquery_checks/Windows/kiro_amazon_q_detection.yml
title: Kiro (Amazon Q) AI Developer Tool Detection
id: 6f7a8b9c0d1e2f3a4b5c6d7e8f9a0b1c
description: |
    Detects the presence of Kiro, AWS's AI-powered developer tool and successor to Amazon Q
    Developer CLI. Kiro is an agentic coding assistant with deep AWS integration that can read
    and modify files, execute shell commands, run tests, and interact directly with AWS services
    and credentials available on the device. Its MCP (Model Context Protocol) support enables
    broad tool integrations including access to databases, APIs, and cloud infrastructure.
    Because Kiro authenticates through AWS IAM Identity Center, a compromised or misconfigured
    Kiro session may have the same AWS permissions as the authenticated user, including access
    to production secrets. This detection also covers residual Amazon Q Developer artifacts
    (config stored in %USERPROFILE%\.amazonq\) that persist during the Q-to-Kiro migration.
    Detection evaluates startup items, file paths, running processes, installed programs, and
    prefetch entries. A score greater than 1 (at least two positive indicators) is required to
    confirm detection.
references:
    - https://kiro.dev
    - https://docs.aws.amazon.com/amazonq/latest/qdeveloper-ug/upgrade-to-kiro.html
    - https://www.okta.com/blog/threat-intelligence/detecting-openclaw-advanced-posture-checks/
author:
    - rafa.bono@okta.com
platform: Windows
query: |
    WITH startup_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM startup_items
        WHERE name LIKE '%kiro%'
            OR name LIKE '%amazon-q%'
    ),
    file_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM file
        WHERE path LIKE '%\Users\%%\AppData\Local\Programs\Kiro\Kiro.exe'
            OR path LIKE '%\Users\%%\.kiro\settings\cli.json'
            OR path LIKE '%\Users\%%\.kiro\settings\mcp.json'
            -- Legacy Amazon Q Developer paths (pre-migration)
            OR path LIKE '%\Users\%%\.amazonq\scopes\scope.json'
            -- AWS Toolkit VS Code extension (ships Kiro/Q integration)
            OR path LIKE '%\Users\%%\.vscode\extensions\amazonwebservices.aws-toolkit-vscode-%'
    ),
    process_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM processes
        WHERE name LIKE '%kiro%'
            OR name LIKE '%amazon-q%'
            OR cmdline LIKE '%kiro-cli%'
    ),
    programs_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM programs
        WHERE name LIKE '%kiro%'
            OR name LIKE '%amazon q%'
    ),
    prefetch_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM prefetch
        WHERE path LIKE '%\kiro%'
            OR path LIKE '%\amazon-q%'
    ),
    sockets_kiro AS (
        SELECT COALESCE(COUNT(*), 0) AS total
        FROM process_open_sockets
        WHERE path LIKE '%\kiro%'
    ),
    final_score AS (
        SELECT
            startup_kiro.total
            + file_kiro.total
            + process_kiro.total
            + programs_kiro.total
            + prefetch_kiro.total
            + sockets_kiro.total
            AS score
        FROM startup_kiro, file_kiro, process_kiro, programs_kiro, prefetch_kiro, sockets_kiro
    )
    SELECT
        CASE
            WHEN score <= 1 THEN 0
            WHEN score > 1 THEN 1
        END AS kiro_detected
    FROM final_score
    ;
